# dotfiles/.config/homebrew/com.nxc.backupbrew.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nxc.backupbrew</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nxc/ghq/github.com/archcorsair/nix-macos/dotfiles/.config/homebrew/backup_brew.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer> <key>Hour</key>
        <integer>14</integer> <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
